# .config/quickshell/ii/distro/arch/ba-shell-git/ba-shell-git.install
post_install() {
  cat <<'EOF'
BlueArchiveOS has been installed.

Recommended next steps:
  - Start the shell manually with: ba run
  - Or enable the user service with:
      systemctl --user enable --now ba.service

Package-managed installs keep user config in:
  ~/.config/illogical-impulse/config.json

The packaged runtime lives in:
  /usr/share/quickshell/ba
EOF
}

post_upgrade() {
  cat <<'EOF'
BlueArchiveOS has been upgraded.

Recommended next steps:
  - Restart the shell with: ba restart
  - Package-managed installs use pacman for updates.
  - Runtime migrations, when needed, are applied by BlueArchiveOS's packaged lifecycle.
EOF
}
